# .config/systemd/fetch-emails.service
[Unit]
Description=Fetch Emails Service

[Service]
Type=simple
ExecStart=/usr/bin/mbsync -a
